# crontab
# Run the app every day at 2:15pm
15 16 * * * root PATH=/usr/local/bin:$PATH /app/run.sh >> /var/log/cron.log 2>&1
